# deploy/atlas.service
# atlas.service — systemd unit to serve an Atlas mind without Docker.
#
# Installation:
#   1. Adjust WorkingDirectory (the ENGINE repo: src/, which bundles the
#      viewer assets under src/web/ and src/templates/), the mind path in
#      ExecStart, and User/Group (dedicated NON-root user, owner of the
#      mind — the server writes and commits to it).
#   2. cp deploy/atlas.service /etc/systemd/system/atlas.service
#   3. systemctl daemon-reload && systemctl enable --now atlas
#
# The mind config (port, auth, trusted_ip_header…) lives in
# <mind>/atlas.toml; secrets can also be passed via Environment= or
# EnvironmentFile= below (e.g. SESSION_SECRET).

[Unit]
Description=Atlas knowledge base server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=atlas
Group=atlas
# The Atlas engine repo (src/ bundles the viewer assets under src/web/ and src/templates/).
WorkingDirectory=/opt/atlas
# Serves the given mind (equivalent to ATLAS_MIND); port via atlas.toml or --port.
ExecStart=/usr/bin/python3 -m src.cli serve /var/lib/atlas/mind
# Secrets outside the mind if needed:
# Environment=SESSION_SECRET=change-me-secret-long-and-random
# EnvironmentFile=/etc/atlas/atlas.env
Restart=on-failure
RestartSec=5
# Give the SIGTERM handler time to push pending writes.
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
